# systemd/autosnap-scheduler.service
[Unit]
Description=autosnap scheduler tick

[Service]
Type=oneshot
WorkingDirectory=/opt/autosnap
ExecStart=/usr/bin/python3 /opt/autosnap/autosnap.py tick
